# services/odoo/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:postgresql]
command=/bin/su postgres -c '/usr/lib/postgresql/9.3/bin/postgres \
  -D /var/lib/postgresql/9.3/main \
  -c config_file=/etc/postgresql/9.3/main/postgresql.conf'

[program:openerp]
command=start-stop-daemon --start --quiet --pidfile /var/run/openerp-server.pid --chuid openerp --make-pidfile --exec /usr/bin/openerp-server -- --config=/etc/openerp/openerp-server.conf --logfile=/var/log/openerp/openerp-server.log
