# openstack/dashboard/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/usr/sbin/apache2 -D FOREGROUND
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:memcached]
command=memcached
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autostart=true
