# bigdata/ambari/ambari-serf/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
 
[program:sshd]
command=/usr/sbin/sshd -D
 
[program:serf]
command=/root/serf/bin/serf_init.sh
