# remote-exec/gunnery/supervisord.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /etc/uwsgi/apps-enabled/gunnery.ini

[program:nginx-app]
command = /usr/sbin/nginx
