# services/graphLab/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:ipython]
command=ipython notebook --no-browser --ip=0.0.0.0 --port 8888/home/node_modules/codebox/bin/codebox.js run /home/myworkspace -u notroot:Root123# -p 7777
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autostart=true